# etc/systemd/fpp-wifi-monitor.service
[Unit]
Description=FPP WiFi Network Monitor
After=network.target hostapd.service
Wants=hostapd.service

[Service]
Type=simple
ExecStart=/opt/fpp/scripts/wifi_monitor.sh
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target